# nginx/openresty
#!/sbin/openrc-run

command="/usr/local/openresty/bin/openresty"
command_args="-g \"daemon off;\""
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"
output_log="/var/log/openresty.log"
