# tika/tika
#!/sbin/openrc-run

command="java"
command_args="-jar -Dlog4j.configurationFile=/app/log4j2.xml /app/tika-server-standard-2.8.0.jar > /var/log/tika-server.log"
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"